# templates/etc/init/openvpn.conf
description "Enstablish open vpn connection"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
chdir /etc/openvpn

exec openvpn France.ovpn
